# docker/dockerfiles/02_redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/opt/bin/redis-server

[program:ssh]
command=/usr/sbin/sshd -D
